# supervisord.d/sender.ini
[program:sender]
command=/home/falcon-ng/falcon-ng-sender

autostart=true
autorestart=true
startsecs=1
startretries=10000000
redirect_stderr=true
directory=/home/falcon-ng
stdout_logfile=/home/falcon-ng/logs/sender.stdout.log
user=root

stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
